# deploy/Typora.ini
[group:Typora]
programs = Typora-8101, Typora-8102

[program:Typora-8101]
directory = /home/www/Typora
command = /home/www/Typora/venv/bin/python3 run.py --port=8101
environment = ENV="prod"
autorestart = true
stdout_logfile = /var/log/supervisor/Typora.log
stderr_logfile = /var/log/supervisor/Typora.err
startsecs = 1

[program:Typora-8102]
directory = /home/www/Typora
command = /home/www/Typora/venv/bin/python3 run.py --port=8102
environment = ENV="prod"
autorestart = true
stdout_logfile = /var/log/supervisor/Typora.log
stderr_logfile = /var/log/supervisor/Typora.err
startsecs = 5