# deploy/typora-blog.ini
[group:typora-blog]
programs = 8101, 8102

[program:8101]
directory = /data/typora-blog
command = /data/typora-blog/venv/bin/python3 run.py --port=8101
environment = ENV="prod"
autorestart = true
stderr_logfile = /data/log/typora-blog.log

[program:8102]
directory = /data/typora-blog
command = /data/typora-blog/venv/bin/python3 run.py --port=8102
environment = ENV="prod"
autorestart = true
stderr_logfile = /data/log/typora-blog.log