# deploy/Typora-Blog.ini
[group:Typora-Blog]
programs = 8101, 8102

[program:8101]
directory = /home/www/Typora-Blog
command = /home/www/Typora-Blog/venv/bin/python3 run.py --port=8101
environment = ENV="prod"
autorestart = true
stdout_logfile = /var/log/supervisor/Typora-Blog.log
stderr_logfile = /var/log/supervisor/Typora-Blog.err

[program:8102]
directory = /home/www/Typora-Blog
command = /home/www/Typora-Blog/venv/bin/python3 run.py --port=8102
environment = ENV="prod"
autorestart = true
stdout_logfile = /var/log/supervisor/Typora-Blog.log
stderr_logfile = /var/log/supervisor/Typora-Blog.err